# upstart-x.conf
start on runlevel [2]
stop on runlevel [!2]

respawn
exec su -l -- marty -c "sh -l -c startx" > /dev/console 2>&1
